# scripts/pingbeacon/.svn/text-base/all.svn-base
#!/bin/sh /etc/rc.common
# IVI start script
# Copyright (C) 2009 Tsinghua University
IVIPATH="/home/jacky/src/ivimod/trunk"
. $IVIPATH/scripts/pingbeacon/startup.sh
. $IVIPATH/scripts/pingbeacon/addr.sh
. $IVIPATH/scripts/pingbeacon/route.sh
. $IVIPATH/scripts/pingbeacon/ivimap.sh
. $IVIPATH/scripts/pingbeacon/log.sh
echo start ivi
